# docs/systemd-with-workers/system/matrix-relapse.target
[Unit]
Description=Relapse parent target
After=network.target

[Install]
WantedBy=multi-user.target
